# X3_SDK_LINUX/X3M_SDK_BR22_20230128-1201/board_support_package/platform_source_code/system/rootfs_yocto/root/etc/udev/rules.d/49-usb.rules
# automounting usb flash drives
KERNEL=="sd[a-z][1-9]", ACTION=="add", RUN+="/bin/mkdir -p /media/%k"
KERNEL=="sd[a-z][1-9]", ACTION=="add", RUN+="/bin/mount /dev/%k /media/%k"

# umount usb flash drivers
KERNEL=="sd[a-z][1-9]", ACTION=="remove", RUN+="/bin/sync"
KERNEL=="sd[a-z][1-9]", ACTION=="remove", RUN+="/bin/umount -l /media/%k"
KERNEL=="sd[a-z][1-9]", ACTION=="remove", RUN+="/bin/rm -rf /media/%k"
