# X3_SDK_LINUX/X3M_SDK_BR22_20230128-1201/board_support_package/platform_source_code/system/rootfs_yocto/root/etc/udev/rules.d/50-adbd.rules
#enable adb when usb switch to device mode
ACTION=="remove", KERNEL=="xhci-hcd.0.auto", RUN+="/etc/init.d/adbd adbd_hotplug"
